# systemd/backupfull.service
[Unit]
Description=Crear backup completo

[Service]
User=root
Type=simple
ExecStart=/ path to backupfull.sh

[Install]
WantedBy=multi-user.target
