# systemd/backupfull.timer
[Unit]
Description=Ejecuta la creacion de backup completo

[Timer]
OnCalendar= por determinar
Persistent=true
 
[Install]
WantedBy=multi-user.target
